# scripts/hotspot/hotspot.service
[Unit]
Description=Start Robotixx_MuSoHu WiFi Hotspot at boot
After=network.target NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/home/jetson/MuSoHu/scripts/hotspot/start-hotspot.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target
