# scripts/musohu-startup.service
[Unit]
Description=MuSoHu Project Startup Service
After=docker.service
Requires=docker.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/jetson/MuSoHu/scripts/start.sh
WorkingDirectory=/home/jetson/MuSoHu
RemainAfterExit=yes
User=jetson
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
